# root/etc/cron.d/save
0 */12 * * * abc /usr/local/bin/terrariad save