# root/etc/services.d/terraria/run
#!/usr/bin/with-contenv bash
echo "Starting Terraria Server."

CMD="/app/TerrariaServer -x64 -config /config/serverconfig.txt -banlist /config/banlist.txt"

exec s6-setuidgid abc screen -DmS terraria ${CMD}